# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-subway/run
#!/command/with-contenv sh

echo "Launching Subway Web Server"
cd /app && bin/puma -C config/puma.rb | tee log/production.log
