# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-subway-worker/run
#!/command/with-contenv sh

echo "Launching Subway Worker Service"
cd /app && bin/good_job start
